# voice_assistant_python.service
[Unit]
Description=Voice Assistant Daemon — Python (Mic → ASR → LLM → TTS → Speaker)
After=network.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu

WorkingDirectory=/opt/voice_assistant_app/voice_assistant_app
ExecStart=/opt/voice_assistant_app/voice_assistant_app/venv/bin/python /opt/voice_assistant_app/voice_assistant_app/voice_assistant.py

# ── Configuration ─────────────────────────────────────────────────────────────
# All settings are loaded from .env in the working directory.
# Copy .env.example → .env on the robot and fill in your values:
#
#   scp .env.example ubuntu@<robot>:/opt/voice_assistant_app/voice_assistant_app/.env
#   ssh ubuntu@<robot> "nano /opt/voice_assistant_app/voice_assistant_app/.env"
#
# Variables set here as Environment= override .env values.
EnvironmentFile=-/opt/voice_assistant_app/voice_assistant_app/.env

# Allow access to audio devices
SupplementaryGroups=audio
# Connect to user's PulseAudio session (needed for mic/speaker access)
Environment=PULSE_RUNTIME_PATH=/run/user/1000/pulse
Environment=XDG_RUNTIME_DIR=/run/user/1000

Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=voice_assistant_python

[Install]
WantedBy=multi-user.target
